# services/news/run
#!/usr/bin/env bash
HOME=/home/chris s6-setuidgid chris /home/chris/.nix-profile/bin/get-news
